# system_files/esfa_checker.timer
[Unit]
Description=Run ESFA Submission Checker every 6 hours

[Timer]
OnUnitActiveSec=10s
Unit=esfa_checker.service

[Install]
WantedBy=timers.target