# system_files/esfa_checker.timer
[Unit]
Description=Run ESFA Submission Checker every 4 hours

[Timer]
OnUnitActiveSec=4h
Unit=esfa_checker.service

[Install]
WantedBy=timers.target